# bots/sigma/launchd/ai.sigma.daily-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.sigma.daily-report</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/tsx</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/sigma/ts/src/sigma-daily-report.ts</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/alexlee</string>
    <key>MODE</key>
    <string>ops</string>
    <key>SIGMA_TELEGRAM_ENHANCED</key>
    <string>false</string>
    <key>SIGMA_HTTP_PORT</key>
    <string>4010</string>
    <key>HUB_BASE_URL</key>
    <string>http://localhost:7788</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/sigma-daily-report.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/sigma-daily-report.err.log</string>

  <!-- 매일 06:30 KST (UTC 21:30 전날) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>21</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
