# bots/sigma/launchd/ai.sigma.observation-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.sigma.observation-review</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/tsx</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/sigma/scripts/runtime-sigma-observation-review.ts</string>
    <string>--telegram</string>
    <string>--record</string>
    <string>--json</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/alexlee</string>
    <key>MODE</key>
    <string>ops</string>
    <key>HUB_BASE_URL</key>
    <string>http://localhost:7788</string>
    <key>SIGMA_OBSERVATION_REVIEW_TELEGRAM</key>
    <string>true</string>
    <key>SIGMA_OBSERVATION_RECORD</key>
    <string>true</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/sigma-observation-review.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/sigma-observation-review.err.log</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
