# bots/sigma/launchd/ai.sigma.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.sigma.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/mix</string>
    <string>shadow</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/sigma/elixir</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/alexlee</string>
    <key>MODE</key>
    <string>ops</string>
    <key>MIX_ENV</key>
    <string>dev</string>
    <key>SIGMA_V2_ENABLED</key>
    <string>true</string>
    <key>SIGMA_TIER2_AUTO_APPLY</key>
    <string>false</string>
    <key>SIGMA_MCP_SERVER_ENABLED</key>
    <string>false</string>
    <key>SIGMA_GEPA_ENABLED</key>
    <string>false</string>
    <key>SIGMA_SELF_RAG_ENABLED</key>
    <string>false</string>
    <key>SIGMA_MAPEK_ENABLED</key>
    <string>false</string>
    <key>SIGMA_SELF_REWARDING_ENABLED</key>
    <string>false</string>
    <key>SIGMA_AGENTIC_RAG_ENABLED</key>
    <string>false</string>
    <key>SIGMA_UCB_ENABLED</key>
    <string>false</string>
    <key>SIGMA_TELEGRAM_ENHANCED</key>
    <string>false</string>
    <key>SIGMA_HTTP_PORT</key>
    <string>4010</string>
    <key>SIGMA_LLM_DAILY_BUDGET_USD</key>
    <string>10.00</string>
    <key>LLM_HUB_ROUTING_ENABLED</key>
    <string>false</string>
    <key>LLM_HUB_ROUTING_SHADOW</key>
    <string>true</string>
    <key>HUB_BASE_URL</key>
    <string>http://localhost:7788</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/sigma-daily.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/sigma-daily.err.log</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
