# bots/worker/launchd/ai.worker.claude-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
	</dict>
	<key>Label</key>
	<string>ai.worker.claude-monitor</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/alexlee/projects/ai-agent-system/node_modules/.bin/tsx</string>
		<string>scripts/claude-api-monitor.ts</string>
		<string>--alert-only</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/Users/alexlee/.openclaw/workspace/logs/claude-monitor-error.log</string>
	<key>StandardOutPath</key>
	<string>/Users/alexlee/.openclaw/workspace/logs/claude-monitor.log</string>
	<key>StartInterval</key>
	<integer>60</integer>
	<key>WorkingDirectory</key>
	<string>/Users/alexlee/projects/ai-agent-system/bots/worker</string>
</dict>
</plist>
