# bots/worker/launchd/ai.worker.health-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/alexlee</string>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
	</dict>
	<key>Label</key>
	<string>ai.worker.health-check</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/alexlee/projects/ai-agent-system/node_modules/.bin/tsx</string>
		<string>/Users/alexlee/projects/ai-agent-system/bots/worker/scripts/health-check.ts</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/tmp/worker-health-check.log</string>
	<key>StandardOutPath</key>
	<string>/tmp/worker-health-check.log</string>
	<key>StartInterval</key>
	<integer>600</integer>
</dict>
</plist>
