# bots/worker/launchd/ai.worker.task-runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.worker.task-runner</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/alexlee/.nvm/versions/node/v24.13.1/bin/node</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/worker/src/task-runner.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/worker</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/alexlee</string>
    <key>PATH</key>
    <string>/Users/alexlee/.nvm/versions/node/v24.13.1/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>

  <key>KeepAlive</key>
  <true/>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/alexlee/.openclaw/workspace/logs/worker-task-runner.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/alexlee/.openclaw/workspace/logs/worker-task-runner-error.log</string>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
