# bots/worker/launchd/ai.worker.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.worker.web</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/worker/scripts/start-worker-web.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/worker</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>WORKER_PORT</key>
    <string>4000</string>
  </dict>

  <key>KeepAlive</key>
  <true/>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/alexlee/.openclaw/workspace/logs/worker-web.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/alexlee/.openclaw/workspace/logs/worker-web-error.log</string>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
