# bots/worker/launchd/ai.worker.nextjs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>CLAUDE_CODE_NAME</key>
		<string>worker-ops</string>
		<key>CLAUDE_CODE_SETTINGS</key>
		<string>/Users/alexlee/.openclaw/.claude/worker-ops.settings.json</string>
		<key>NEXT_PUBLIC_API_URL</key>
		<string>http://localhost:4000</string>
		<key>NODE_ENV</key>
		<string>production</string>
		<key>OPENCLAW_AGENT</key>
		<string>worker-ops</string>
		<key>PATH</key>
		<string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>ai.worker.nextjs</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/npm</string>
		<string>run</string>
		<string>start</string>
		<string>--</string>
		<string>-H</string>
		<string>0.0.0.0</string>
		<string>-p</string>
		<string>4001</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/Users/alexlee/.openclaw/workspace/logs/worker-nextjs-error.log</string>
	<key>StandardOutPath</key>
	<string>/Users/alexlee/.openclaw/workspace/logs/worker-nextjs.log</string>
	<key>ThrottleInterval</key>
	<integer>10</integer>
	<key>WorkingDirectory</key>
	<string>/Users/alexlee/projects/ai-agent-system/bots/worker/web</string>
</dict>
</plist>
