# bots/darwin/launchd/ai.darwin.weekly.autonomous.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
                       "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.darwin.weekly.autonomous</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/alexlee/projects/ai-agent-system/node_modules/.bin/tsx</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/darwin/lib/research-scanner.ts</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key><string>/Users/alexlee</string>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>MIX_ENV</key><string>prod</string>
    <key>MODE</key><string>ops</string>
    <key>DARWIN_AUTONOMY_LEVEL</key><string>5</string>
    <key>DARWIN_V2_ENABLED</key><string>true</string>
    <key>DARWIN_SHADOW_MODE</key><string>false</string>
    <key>DARWIN_KILL_SWITCH</key><string>false</string>
    <key>DARWIN_L5_ENABLED</key><string>true</string>
    <key>DARWIN_HTTP_PORT</key><string>4020</string>
    <key>DARWIN_LLM_DAILY_BUDGET_USD</key><string>10.00</string>
    <key>DARWIN_TIER2_AUTO_APPLY</key><string>true</string>
    <key>DARWIN_MCP_SERVER_ENABLED</key><string>false</string>
    <key>DARWIN_ESPL_ENABLED</key><string>false</string>
    <key>DARWIN_SELF_RAG_ENABLED</key><string>false</string>
    <key>DARWIN_PRINCIPLE_SEMANTIC_CHECK_ENABLED</key><string>false</string>
    <!-- Phase A: 9팀 기술 요청 통합 (Kill Switch 기본 false) -->
    <key>DARWIN_TEAM_INTEGRATION_ENABLED</key><string>false</string>
    <!-- Phase B: Hypothesis Engine (Sakana AI Scientist 패턴, Kill Switch 기본 false) -->
    <key>DARWIN_HYPOTHESIS_ENGINE_ENABLED</key><string>false</string>
    <key>DARWIN_HYPOTHESIS_LLM_DAILY_BUDGET_USD</key><string>2.0</string>
    <!-- Phase C: MEASURE Stage 24h/7d/30d 효과 측정 (Kill Switch 기본 false) -->
    <key>DARWIN_MEASURE_STAGE_ENABLED</key><string>false</string>
    <key>LLM_HUB_ROUTING_ENABLED</key><string>true</string>
    <key>LLM_HUB_ROUTING_SHADOW</key><string>false</string>
    <key>HUB_BASE_URL</key><string>http://localhost:7788</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>0</integer>
    <key>Hour</key><integer>5</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key><string>/tmp/darwin-weekly-autonomous.log</string>
  <key>StandardErrorPath</key><string>/tmp/darwin-weekly-autonomous.err.log</string>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
