# bots/edu-x/launchd/ai.edux.daily-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/alexlee</string>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>EDUX_DIGEST_DRY_RUN</key>
		<string>false</string>
		<key>EDUX_DIGEST_TELEGRAM_TARGET</key>
		<string>live</string>
		<key>EDUX_DIGEST_LINK_PREVIEW</key>
		<string>false</string>
		<key>EDUX_DIGEST_SEND_THUMBNAIL</key>
		<string>false</string>
	</dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>ai.edux.daily-digest</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/node</string>
		<string>--env-file=/Users/alexlee/projects/ai-agent-system/.env</string>
		<string>/Users/alexlee/projects/ai-agent-system/bots/edu-x/scripts/runtime-edux-daily-digest.ts</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/Users/alexlee/projects/ai-agent-system/bots/edu-x/daily-digest.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/alexlee/projects/ai-agent-system/bots/edu-x/daily-digest.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>10</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>TimeOut</key>
	<integer>600</integer>
</dict>
</plist>
