# bots/claude/launchd/ai.claude.dexter.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.claude.dexter.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/claude/src/dexter.js</string>
    <string>--daily-report</string>
    <string>--telegram</string>
  </array>

  <!-- 매일 오전 8시 KST -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/claude</string>

  <key>StandardOutPath</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/claude/dexter-daily.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/claude/dexter-daily.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <key>HOME</key>
    <string>/Users/alexlee</string>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
