# bots/claude/launchd/ai.claude.archer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/alexlee</string>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
		<key>OPENCLAW_AGENT</key>
		<string>claude-ops</string>
		<key>CLAUDE_CODE_NAME</key>
		<string>claude-ops</string>
		<key>CLAUDE_CODE_SETTINGS</key>
		<string>/Users/alexlee/.openclaw/.claude/claude-ops.settings.json</string>
	</dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>ai.claude.archer</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/node</string>
		<string>/Users/alexlee/projects/ai-agent-system/bots/claude/src/archer.js</string>
		<string>--telegram</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/Users/alexlee/projects/ai-agent-system/bots/claude/archer.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/alexlee/projects/ai-agent-system/bots/claude/archer.log</string>
	<key>StartCalendarInterval</key>
	<array>
		<dict>
			<key>Hour</key>
			<integer>9</integer>
			<key>Minute</key>
			<integer>0</integer>
			<key>Weekday</key>
			<integer>1</integer>
		</dict>
		<dict>
			<key>Hour</key>
			<integer>9</integer>
			<key>Minute</key>
			<integer>0</integer>
			<key>Weekday</key>
			<integer>3</integer>
		</dict>
		<dict>
			<key>Hour</key>
			<integer>9</integer>
			<key>Minute</key>
			<integer>0</integer>
			<key>Weekday</key>
			<integer>5</integer>
		</dict>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/alexlee/projects/ai-agent-system/bots/claude</string>
</dict>
</plist>
