# bots/claude/launchd/ai.claude.refactor-cycle.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/alexlee</string>
		<key>MODE</key>
		<string>ops</string>
		<key>NODE_ENV</key>
		<string>production</string>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>REFACTORER_CYCLE_MODE</key>
		<string>active</string>
		<key>REFACTORER_APPLY_ENABLED</key>
		<string>true</string>
		<key>REFACTORER_DIRTY_SCOPE</key>
		<string>workspace</string>
		<key>REFACTORER_AUTOFIX_ENABLED</key>
		<string>true</string>
		<key>REFACTORER_AUTOFIX_MAX_ATTEMPTS</key>
		<string>2</string>
		<key>REFACTORER_AUTOFIX_ESTIMATED_COST_LIMIT</key>
		<string>1.0</string>
		<key>REFACTORER_APPLY_PUSH</key>
		<string>true</string>
		<key>REFACTORER_APPLY_STRICT_GATE</key>
		<string>true</string>
		<key>REFACTORER_APPLY_MAX_PER_CYCLE</key>
		<string>3</string>
	</dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>ai.claude.refactor-cycle</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/node</string>
		<string>/Users/alexlee/projects/ai-agent-system/bots/claude/scripts/refactor-cycle-runner.ts</string>
		<string>--json</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/Users/alexlee/.ai-agent-system/logs/claude-refactor-cycle-error.log</string>
	<key>StandardOutPath</key>
	<string>/Users/alexlee/.ai-agent-system/logs/claude-refactor-cycle.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>3</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>WorkingDirectory</key>
	<string>/Users/alexlee/projects/ai-agent-system</string>
</dict>
</plist>
