# bots/investment/launchd/ai.investment.crypto.validation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.investment.crypto.validation</string>
  <!-- Phase 3-B: 암호화폐 검증거래 레일 (validation mode) -->
  <!-- normal 레일과 다른 guard scope / 로그 경로를 사용 -->
  <!-- 설치: cp launchd/ai.investment.crypto.validation.plist ~/Library/LaunchAgents/ -->
  <!--        launchctl load ~/Library/LaunchAgents/ai.investment.crypto.validation.plist -->
  <!-- 주의: validation 전용 sizing/risk 정책을 붙이기 전에는 운영 판단 하에만 활성화 -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/alexlee/.nvm/versions/node/v24.13.1/bin/node</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/investment/markets/crypto.js</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/investment</string>
  <!-- 15분 주기 (validation canary) -->
  <key>StartInterval</key>
  <integer>900</integer>
  <key>StandardOutPath</key>
  <string>/tmp/investment-crypto-validation.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/investment-crypto-validation.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/alexlee</string>
    <key>PATH</key>
    <string>/Users/alexlee/.nvm/versions/node/v24.13.1/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>PAPER_MODE</key>
    <string>false</string>
    <key>INVESTMENT_TRADE_MODE</key>
    <string>validation</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
