# bots/investment/launchd/ai.investment.crypto.validation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/alexlee</string>
		<key>INVESTMENT_TRADE_MODE</key>
		<string>validation</string>
		<key>NODE_ENV</key>
		<string>production</string>
		<key>PAPER_MODE</key>
		<string>false</string>
		<key>PATH</key>
		<string>/Users/alexlee/.nvm/versions/node/v24.13.1/bin:/usr/local/bin:/usr/bin:/bin</string>
	</dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>ai.investment.crypto.validation</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/alexlee/.nvm/versions/node/v24.13.1/bin/node</string>
		<string>/Users/alexlee/projects/ai-agent-system/bots/investment/markets/crypto.js</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/tmp/investment-crypto-validation.err.log</string>
	<key>StandardOutPath</key>
	<string>/tmp/investment-crypto-validation.log</string>
	<key>StartInterval</key>
	<integer>900</integer>
	<key>WorkingDirectory</key>
	<string>/Users/alexlee/projects/ai-agent-system/bots/investment</string>
</dict>
</plist>
