# bots/investment/launchd/ai.investment.commander.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>LLM_USE_OAUTH_PRIMARY</key>
		<string>true</string>
		<key>LLM_TEAM_SELECTOR_VERSION</key>
		<string>v3.0_oauth_4</string>
		<key>LLM_TEAM_SELECTOR_AB_PERCENT</key>
		<string>100</string>
		<key>LLM_TEAM_SELECTOR_VERSION_PCT</key>
		<string>100</string>
		<key>INVESTMENT_LLM_HUB_ENABLED</key>
		<string>true</string>
		<key>LUNA_AGENT_LLM_ROUTING_ENABLED</key>
		<string>true</string>
		<key>INVESTMENT_LLM_DIRECT_FALLBACK</key>
		<string>false</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>ai.investment.commander</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/node</string>
		<string>/Users/alexlee/projects/ai-agent-system/bots/investment/luna-commander.cjs</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/alexlee/.ai-agent-system/logs/luna-commander-error.log</string>
	<key>StandardOutPath</key>
	<string>/Users/alexlee/.ai-agent-system/logs/luna-commander.log</string>
	<key>ThrottleInterval</key>
	<integer>10</integer>
	<key>WorkingDirectory</key>
	<string>/Users/alexlee/projects/ai-agent-system/bots/investment</string>
</dict>
</plist>
