# bots/investment/launchd/ai.investment.unrealized-pnl.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.investment.unrealized-pnl</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/alexlee/projects/ai-agent-system/bots/investment/scripts/update-unrealized-pnl.js</string>
    </array>
    <key>StartInterval</key>
    <integer>600</integer>
    <key>StandardOutPath</key>
    <string>/tmp/unrealized-pnl.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/unrealized-pnl.err.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/alexlee/projects/ai-agent-system</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/alexlee</string>
    </dict>
</dict>
</plist>
