# bots/investment/launchd/ai.investment.runtime-autopilot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/alexlee</string>
		<key>NODE_ENV</key>
		<string>production</string>
		<key>HEPHAESTOS_HOT_PATH_PREFETCH_ENABLED</key>
		<string>1</string>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>PAPER_MODE</key>
		<string>false</string>
		<key>LUNA_DELEGATED_AUTHORITY_ENABLED</key>
		<string>true</string>
		<key>LUNA_MASTER_REPORT_ONLY</key>
		<string>true</string>
		<key>LUNA_DELEGATED_AUTHORITY_REQUIRE_FINAL_GATE</key>
		<string>true</string>
		<key>LUNA_DELEGATED_MAX_TRADE_USDT</key>
		<string>50</string>
		<key>LUNA_DELEGATED_MAX_DAILY_USDT</key>
		<string>200</string>
		<key>LUNA_DELEGATED_MAX_OPEN_POSITIONS</key>
		<string>2</string>
	</dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>ai.investment.runtime-autopilot</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/node</string>
		<string>/Users/alexlee/projects/ai-agent-system/bots/investment/scripts/runtime-position-runtime-autopilot.ts</string>
		<string>--execute</string>
		<string>--apply-tuning</string>
		<string>--execute-dispatch</string>
		<string>--confirm=position-runtime-autopilot</string>
		<string>--json</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/tmp/investment-runtime-autopilot.err.log</string>
	<key>StandardOutPath</key>
	<string>/tmp/investment-runtime-autopilot.log</string>
	<key>StartInterval</key>
	<integer>120</integer>
	<key>WorkingDirectory</key>
	<string>/Users/alexlee/projects/ai-agent-system/bots/investment</string>
</dict>
</plist>
