# bots/investment/launchd/ai.investment.market-alert-crypto-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.investment.market-alert-crypto-daily</string>

  <!-- 암호화폐 일일 요약 보고 — KST 09:00 (UTC 00:00) -->
  <!-- 설치: cp bots/investment/launchd/ai.investment.market-alert-crypto-daily.plist ~/Library/LaunchAgents/ -->
  <!--        launchctl load ~/Library/LaunchAgents/ai.investment.market-alert-crypto-daily.plist -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/alexlee/.nvm/versions/node/v24.13.1/bin/node</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/investment/scripts/market-alert.js</string>
    <string>--market=crypto</string>
    <string>--event=daily</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/investment</string>

  <!-- KST 09:00 = UTC 00:00 (국내장 시작 알림과 동시 실행) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>0</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/investment-market-alert-crypto-daily.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/investment-market-alert-crypto-daily.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/alexlee</string>
    <key>PATH</key>
    <string>/Users/alexlee/.nvm/versions/node/v24.13.1/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>PAPER_MODE</key>
    <string>false</string>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
