# bots/investment/launchd/ai.luna.fx-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.luna.fx-refresh</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/psql</string>
    <string>-d</string>
    <string>jay</string>
    <string>-c</string>
    <string>REFRESH MATERIALIZED VIEW CONCURRENTLY investment.v_trades_real_usd;</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <!-- 매일 06:00 KST = 21:00 UTC (launchd는 KST 기준, UTC 변환 금지) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/tmp/ai.luna.fx-refresh.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/ai.luna.fx-refresh.err.log</string>
</dict>
</plist>
