# bots/investment/launchd/ai.luna.7day-natural-checkpoint.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.luna.7day-natural-checkpoint</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/investment/scripts/runtime-luna-7day-natural-checkpoint.ts</string>
    <string>--write</string>
    <string>--json</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/investment</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>LUNA_7DAY_NATURAL_CHECKPOINT_ENABLED</key>
    <string>true</string>
  </dict>
  <key>StartInterval</key>
  <integer>86400</integer>
  <key>StandardOutPath</key>
  <string>/tmp/ai.luna.7day-natural-checkpoint.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/ai.luna.7day-natural-checkpoint.err.log</string>
</dict>
</plist>
