# bots/investment/launchd/ai.luna.log-rotate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.luna.log-rotate</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/alexlee</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin</string>
    <key>LUNA_LOG_ROTATE_MIN_BYTES</key>
    <string>209715200</string>
    <key>LUNA_LOG_ROTATE_KEEP_LINES</key>
    <string>5000</string>
    <key>LUNA_LOG_ROTATE_KEEP_DAYS</key>
    <string>7</string>
  </dict>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>--disable-warning=DEP0205</string>
    <string>--import</string>
    <string>tsx</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/investment/scripts/luna-log-rotate.ts</string>
    <string>--json</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>20</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/ai.luna.log-rotate.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/ai.luna.log-rotate.err.log</string>
</dict>
</plist>
