# bots/investment/launchd/ai.luna.ppo-retrain-weekly-sun-0200.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.luna.ppo-retrain-weekly-sun-0200</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>--disable-warning=DEP0205</string>
    <string>--import</string>
    <string>tsx</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/investment/scripts/runtime-luna-ppo-retrain-weekly.ts</string>
    <string>--json</string>
    <string>--train</string>
    <string>--confirm=luna-ppo-shadow-train</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/investment</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>LUNA_PPO_RETRAIN_ENABLED</key>
    <string>false</string>
  </dict>
  <!-- 라벨명에 0200 표기 있으나 실제 실행은 일요일 10:00 KST (메모리 분산) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>0</integer>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/tmp/ai.luna.ppo-retrain-weekly-sun-0200.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/ai.luna.ppo-retrain-weekly-sun-0200.err.log</string>
</dict>
</plist>
