# bots/investment/launchd/ai.luna.cache-cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.luna.cache-cleanup</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/alexlee</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin</string>
    <key>LUNA_CACHE_CLEANUP_DISK_THRESHOLD_PCT</key>
    <string>92</string>
  </dict>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>--disable-warning=DEP0205</string>
    <string>--import</string>
    <string>tsx</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/investment/scripts/luna-cache-cleanup.ts</string>
    <string>--json</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>0</integer>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/ai.luna.cache-cleanup.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/ai.luna.cache-cleanup.err.log</string>
</dict>
</plist>
