# bots/investment/launchd/ai.luna.market-gate-30min.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.luna.market-gate-30min</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>--disable-warning=DEP0205</string>
    <string>--import</string>
    <string>tsx</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/investment/scripts/runtime-luna-market-gate.ts</string>
    <string>--json</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>LUNA_ENTRY_TRIGGER_SHADOW</key>
    <string>true</string>
    <key>LUNA_REGIME_EXPANSION_SHADOW</key>
    <string>true</string>
  </dict>
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>RunAtLoad</key>
  <false/>
  <key>KeepAlive</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/tmp/logs/luna-market-gate/out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/logs/luna-market-gate/err.log</string>
  <key>TimeOut</key>
  <integer>600</integer>
</dict>
</plist>
