# bots/investment/launchd/ai.luna.candidate-backtest-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.luna.candidate-backtest-refresh</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/investment/scripts/runtime-luna-candidate-backtest-refresh.ts</string>
    <string>--periods=180</string>
    <string>--json</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/investment</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>LUNA_CANDIDATE_BACKTEST_SHADOW_MODE</key>
    <string>true</string>
    <key>LUNA_BACKTEST_STALE_HOURS</key>
    <string>24</string>
    <key>LUNA_BT_WALK_FORWARD_ENABLED</key>
    <string>true</string>
    <key>LUNA_BT_ROBUST_SELECTION_ENABLED</key>
    <string>true</string>
    <key>LUNA_BT_WF_FOLDS</key>
    <string>5</string>
    <key>LUNA_BT_WF_TRAIN_DAYS</key>
    <string>45</string>
    <key>LUNA_BT_WF_TEST_DAYS</key>
    <string>30</string>
    <key>LUNA_VECTORBT_TIMEOUT_MS</key>
    <string>90000</string>
    <key>LUNA_DSR_GATE_ENABLED</key>
    <string>false</string>
    <key>LUNA_PSR_GATE_ENABLED</key>
    <string>true</string>
    <key>LUNA_PSR_MIN</key>
    <string>0.5</string>
    <key>LUNA_BINANCE_TOP_VOLUME_LIMIT</key>
    <string>50</string>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>KeepAlive</key>
  <false/>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>12</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>TimeOut</key>
  <integer>1800</integer>
  <key>StandardOutPath</key>
  <string>/tmp/logs/luna-candidate-backtest/out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/logs/luna-candidate-backtest/err.log</string>
</dict>
</plist>
