# bots/investment/launchd/ai.luna.hybrid-promotion-gate-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.luna.hybrid-promotion-gate-daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/investment/scripts/luna-shadow-auto-promote.ts</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/investment</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>NODE_ENV</key>
    <string>production</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/alexlee/.ai-agent-system/logs/luna-hybrid-promotion-gate-daily.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/alexlee/.ai-agent-system/logs/luna-hybrid-promotion-gate-daily-error.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
