# bots/investment/launchd/ai.luna.opendart-financial-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.luna.opendart-financial-refresh</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/investment/scripts/runtime-luna-opendart-financial-batch-refresh.ts</string>
    <string>--json</string>
    <string>--bulk</string>
    <string>--network</string>
    <string>--limit=250</string>
    <string>--bsns-years=2024,2023,2022,2021,2020,2019</string>
    <string>--reprt-codes=11011</string>
    <string>--skip-fresh</string>
    <string>--write</string>
    <string>--confirm=luna-opendart-financial-batch-write</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/investment</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>LUNA_KOREA_DATA_SHADOW_MODE</key>
    <string>true</string>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>KeepAlive</key>
  <false/>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>18</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>TimeOut</key>
  <integer>1800</integer>
  <key>StandardOutPath</key>
  <string>/tmp/logs/luna-opendart-financial/out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/logs/luna-opendart-financial/err.log</string>
</dict>
</plist>
