# bots/investment/launchd/ai.luna.weight-adaptive-tuner-daily-0700.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.luna.weight-adaptive-tuner-daily-0700</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/investment/scripts/run-regime-weight-learner.ts</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/investment</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>NODE_ENV</key>
    <string>production</string>
  </dict>
  <!-- 매일 07:00 KST — 전날 거래 분석 + 체제별 fusion/signal 가중치 조정 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/tmp/ai.luna.weight-adaptive-tuner-daily-0700.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/ai.luna.weight-adaptive-tuner-daily-0700.err.log</string>
</dict>
</plist>
