# bots/investment/launchd/ai.luna.universe-refresh-daily-0830.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.luna.universe-refresh-daily-0830</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/investment/scripts/run-universe-refresh.ts</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/investment</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>LUNA_BINANCE_TOP_VOLUME_LIMIT</key>
    <string>50</string>
  </dict>
  <!-- 매일 08:30 KST — 장 시작 전 유니버스 갱신 (국내 Top50 + 해외 Top50 + 크립토 Top50) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/tmp/ai.luna.universe-refresh-daily-0830.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/ai.luna.universe-refresh-daily-0830.err.log</string>
</dict>
</plist>
