# bots/investment/launchd/ai.luna.agent-evolution-weekly-sun-0600.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.luna.agent-evolution-weekly-sun-0600</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/investment/scripts/runtime-luna-agent-evolution.ts</string>
    <string>--json</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/investment</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>NODE_ENV</key>
    <string>production</string>
  </dict>
  <!-- 매주 일요일 06:55 KST — 스케줄 분산(06:00 피크 완화) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>0</integer>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>55</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/tmp/ai.luna.agent-evolution-weekly-sun-0600.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/ai.luna.agent-evolution-weekly-sun-0600.err.log</string>
</dict>
</plist>
