# bots/investment/launchd/ai.luna.crypto-holding-monitor-6h.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.luna.crypto-holding-monitor-6h</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/alexlee/projects/ai-agent-system/dist/daemons/ai.luna.crypto-holding-monitor-6h.mjs</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/alexlee</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>NODE_ENV</key>
    <string>production</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>StartInterval</key>
  <integer>21600</integer>
  <key>TimeOut</key>
  <integer>600</integer>
  <key>StandardOutPath</key>
  <string>/tmp/ai.luna.crypto-holding-monitor-6h.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/ai.luna.crypto-holding-monitor-6h.err.log</string>
</dict>
</plist>
