# scripts/launchd/ai.agent.auto-commit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/alexlee</string>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
	</dict>
	<key>Label</key>
	<string>ai.agent.auto-commit</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/alexlee/projects/ai-agent-system/scripts/auto-commit.sh</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/Users/alexlee/.ai-agent-system/logs/auto-commit.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/alexlee/.ai-agent-system/logs/auto-commit.log</string>
	<key>StartInterval</key>
	<integer>21600</integer>
</dict>
</plist>
