# scripts/launchd/ai.agent.post-reboot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict/>
	<key>Label</key>
	<string>ai.agent.post-reboot</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/alexlee/projects/ai-agent-system/scripts/post-reboot.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/tmp/post-reboot.err.log</string>
	<key>StandardOutPath</key>
	<string>/tmp/post-reboot.log</string>
	<key>WorkingDirectory</key>
	<string>/Users/alexlee/projects/ai-agent-system</string>
</dict>
</plist>
