# scripts/launchd/ai.agent.post-reboot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.agent.post-reboot</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/alexlee/projects/ai-agent-system/scripts/post-reboot.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system</string>

  <!-- 재부팅 후 1회 자동 실행 (KeepAlive 없음) -->
  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/post-reboot.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/post-reboot.err.log</string>
</dict>
</plist>
