# bots/blog/launchd/ai.blog.comfyui.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>COMFYUI_HOME</key>
		<string>/Users/alexlee/ComfyUI</string>
		<key>COMFYUI_HOST</key>
		<string>127.0.0.1</string>
		<key>COMFYUI_MODE</key>
		<string>auto</string>
		<key>COMFYUI_PORT</key>
		<string>8188</string>
		<key>HOME</key>
		<string>/Users/alexlee</string>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>ai.blog.comfyui</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/alexlee/projects/ai-agent-system/scripts/start-comfyui.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/alexlee/projects/ai-agent-system/bots/blog/blog-comfyui.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/alexlee/projects/ai-agent-system/bots/blog/blog-comfyui.log</string>
	<key>WorkingDirectory</key>
	<string>/Users/alexlee/projects/ai-agent-system</string>
</dict>
</plist>
