# bots/blog/launchd/ai.blog.node-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>BLOG_API_PORT</key>
		<string>3100</string>
		<key>BLOG_IMAGE_BASE_URL</key>
		<string>http://127.0.0.1:7860</string>
		<key>BLOG_IMAGE_PROVIDER</key>
		<string>drawthings</string>
		<key>HOME</key>
		<string>/Users/alexlee</string>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>ai.blog.node-server</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/node</string>
		<string>/Users/alexlee/projects/ai-agent-system/bots/blog/api/node-server.ts</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/alexlee/projects/ai-agent-system/bots/blog/blog-node-server.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/alexlee/projects/ai-agent-system/bots/blog/blog-node-server.log</string>
	<key>WorkingDirectory</key>
	<string>/Users/alexlee/projects/ai-agent-system</string>
</dict>
</plist>
