# bots/blog/launchd/ai.blog.neighbor-sympathy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.blog.neighbor-sympathy</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/blog/scripts/run-neighbor-sympathy.js</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>24</integer></dict>
    <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>36</integer></dict>
    <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>12</integer></dict>
    <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>48</integer></dict>
    <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>24</integer></dict>
    <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>36</integer></dict>
    <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>12</integer></dict>
    <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>48</integer></dict>
    <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>24</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>36</integer></dict>
    <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>12</integer></dict>
    <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>48</integer></dict>
    <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>24</integer></dict>
    <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>36</integer></dict>
    <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>12</integer></dict>
    <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>48</integer></dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>MODE</key>
    <string>ops</string>
    <key>BLOG_COMMENTER_TRACE</key>
    <string>false</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/blog/neighbor-sympathy.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/blog/neighbor-sympathy.err.log</string>
</dict>
</plist>
