# bots/blog/launchd/ai.blog.marketing-snapshot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.blog.marketing-snapshot</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/alexlee/projects/ai-agent-system/elixir/team_jay &amp;&amp; mix blog.marketing.snapshot</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system/elixir/team_jay</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/blog/marketing-snapshot.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/blog/marketing-snapshot.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
