# bots/reservation/launchd/ai.ska.ops-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/alexlee</string>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>SKA_OPS_MCP_PORT</key>
		<string>4098</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>ai.ska.ops-mcp</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/node</string>
		<string>--disable-warning=DEP0205</string>
		<string>--import</string>
		<string>tsx</string>
		<string>/Users/alexlee/projects/ai-agent-system/bots/reservation/mcp/ska-ops-mcp/src/server.ts</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/alexlee/projects/ai-agent-system</string>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/tmp/ska-ops-mcp.log</string>
	<key>StandardOutPath</key>
	<string>/tmp/ska-ops-mcp.log</string>
</dict>
</plist>
