# bots/reservation/launchd/ai.ska.forecast-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.ska.forecast-weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/alexlee/projects/ai-agent-system/bots/ska/scripts/run-forecast.sh</string>
        <string>weekly</string>
    </array>

    <!-- 매주 금요일 18:00 — 다음 주 예상 매출 예측 + 텔레그램 전송 -->
    <!-- launchd Weekday: 0=일 1=월 2=화 3=수 4=목 5=금 6=토 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>5</integer>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/ska-forecast-weekly.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ska-forecast-weekly.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/alexlee</string>
        <key>PATH</key>
        <string>/Users/alexlee/projects/ai-agent-system/bots/ska/venv/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        </dict>
</dict>
</plist>
