# bots/reservation/launchd/ai.ska.forecast-monthly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.ska.forecast-monthly</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/ska/scripts/run-forecast-monthly.sh</string>
  </array>

  <!-- 매월 28일 18:00 — 월간 30일 예측 + LLM 모델 진단 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Day</key>
    <integer>28</integer>
    <key>Hour</key>
    <integer>18</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/alexlee</string>
    <key>PATH</key>
    <string>/Users/alexlee/projects/ai-agent-system/bots/ska/venv/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>

  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/ska-forecast-monthly.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/ska-forecast-monthly.err</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
