# bots/reservation/launchd/ai.ska.log-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>CLAUDE_CODE_NAME</key>
		<string>ska-ops</string>
		<key>CLAUDE_CODE_SETTINGS</key>
		<string>/Users/alexlee/.openclaw/.claude/ska-ops.settings.json</string>
		<key>HOME</key>
		<string>/Users/alexlee</string>
		<key>OPENCLAW_AGENT</key>
		<string>ska-ops</string>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
	</dict>
	<key>Label</key>
	<string>ai.ska.log-report</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/alexlee/projects/ai-agent-system/bots/reservation/manual/reports/log-report.sh</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/Users/alexlee/.openclaw/workspace/log-report-err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/alexlee/.openclaw/workspace/log-report-out.log</string>
	<key>StartInterval</key>
	<integer>10800</integer>
</dict>
</plist>
