# bots/reservation/launchd/ai.ska.rebecca-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>CLAUDE_CODE_NAME</key>
		<string>ska-ops</string>
		<key>CLAUDE_CODE_SETTINGS</key>
		<string>/Users/alexlee/.openclaw/.claude/ska-ops.settings.json</string>
		<key>HOME</key>
		<string>/Users/alexlee</string>
		<key>OPENCLAW_AGENT</key>
		<string>ska-ops</string>
		<key>PATH</key>
		<string>/Users/alexlee/projects/ai-agent-system/bots/ska/venv/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
	</dict>
	<key>Label</key>
	<string>ai.ska.rebecca-weekly</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/alexlee/projects/ai-agent-system/bots/ska/scripts/run-rebecca.sh</string>
		<string>weekly</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/tmp/ska-rebecca-weekly.err</string>
	<key>StandardOutPath</key>
	<string>/tmp/ska-rebecca-weekly.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>8</integer>
		<key>Minute</key>
		<integer>5</integer>
		<key>Weekday</key>
		<integer>1</integer>
	</dict>
</dict>
</plist>
