# bots/orchestrator/launchd/ai.orchestrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/alexlee</string>
		<key>MODE</key>
		<string>ops</string>
		<key>HUB_ALARM_USE_CLASS_TOPICS</key>
		<string>true</string>
		<key>NODE_ENV</key>
		<string>production</string>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>PG_DIRECT</key>
		<string>true</string>
		<key>JAY_LLM_DAILY_BUDGET_USD</key>
		<string>5.0</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>ai.orchestrator</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/tsx</string>
		<string>/Users/alexlee/projects/ai-agent-system/bots/orchestrator/src/orchestrator.ts</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/alexlee/.ai-agent-system/logs/orchestrator-error.log</string>
	<key>StandardOutPath</key>
	<string>/Users/alexlee/.ai-agent-system/logs/orchestrator.log</string>
	<key>ThrottleInterval</key>
	<integer>10</integer>
	<key>WorkingDirectory</key>
	<string>/Users/alexlee/projects/ai-agent-system/bots/orchestrator</string>
</dict>
</plist>
