# bots/orchestrator/launchd/ai.orchestrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.orchestrator</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/alexlee/.nvm/versions/node/v24.13.1/bin/node</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/orchestrator/src/mainbot.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/orchestrator</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/alexlee</string>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>PATH</key>
    <string>/Users/alexlee/.nvm/versions/node/v24.13.1/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <!-- 크래시 자동 재시작 -->
  <key>KeepAlive</key>
  <true/>

  <!-- 연속 크래시 방지 (10초 대기) -->
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/alexlee/.openclaw/logs/mainbot.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/alexlee/.openclaw/logs/mainbot-error.log</string>

  <!-- 시스템 부팅 또는 로그인 시 자동 시작 -->
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
