# bots/hub/launchd/ai.hub.ops-console-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.hub.ops-console-tunnel</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/hub/scripts/ops-console-tunnel.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system</string>
  <key>RunAtLoad</key>
  <false/>
  <key>KeepAlive</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/alexlee/.ai-agent-system/logs/ai.hub.ops-console-tunnel.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/alexlee/.ai-agent-system/logs/ai.hub.ops-console-tunnel.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
