# bots/hub/launchd/ai.hub.llm-auto-routing-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.hub.llm-auto-routing-monitor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/tsx</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/hub/scripts/llm-auto-routing-stats.ts</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/alexlee</string>
    <key>MODE</key>
    <string>ops</string>
    <key>LLM_AUTO_ROUTING_ENABLED</key>
    <string>shadow</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/hub-llm-auto-routing-monitor.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/hub-llm-auto-routing-monitor.err.log</string>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Minute</key>
      <integer>10</integer>
    </dict>
    <dict>
      <key>Minute</key>
      <integer>20</integer>
    </dict>
    <dict>
      <key>Minute</key>
      <integer>30</integer>
    </dict>
    <dict>
      <key>Minute</key>
      <integer>40</integer>
    </dict>
    <dict>
      <key>Minute</key>
      <integer>50</integer>
    </dict>
  </array>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
