# bots/hub/launchd/ai.hub.severity-decay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.hub.severity-decay</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/tsx</string>
    <string>/Users/alexlee/projects/ai-agent-system/bots/hub/scripts/severity-decay-runner.ts</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/alexlee</string>
    <key>MODE</key>
    <string>ops</string>
    <key>HUB_ALARM_USE_CLASS_TOPICS</key>
    <string>true</string>
    <key>HUB_SEVERITY_DECAY_ENABLED</key>
    <string>true</string>
    <key>HUB_SEVERITY_DECAY_CRITICAL_HOURS</key>
    <string>24</string>
    <key>HUB_SEVERITY_DECAY_ERROR_DAYS</key>
    <string>7</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/hub-severity-decay.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/hub-severity-decay.err.log</string>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
