# bots/hub/launchd/ai.hub.resource-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.hub.resource-api</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/alexlee/projects/ai-agent-system/node_modules/.bin/tsx</string>
    <string>bots/hub/src/hub.ts</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>MODE</key>
    <string>ops</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HUB_PG_USER</key>
    <string>hub_readonly</string>
    <key>HUB_PG_PASSWORD</key>
    <string>__SET_IN_LOCAL_LAUNCHAGENT__</string>
    <key>HUB_PG_DATABASE</key>
    <string>jay</string>
    <key>HUB_CONTROL_CALLBACK_SECRET</key>
    <string>__SET_IN_LOCAL_LAUNCHAGENT__</string>
    <key>HUB_ALARM_USE_CLASS_TOPICS</key>
    <string>true</string>
    <key>HUB_ALARM_DISPATCH_MODE</key>
    <string>shadow</string>
    <key>HUB_ALARM_LLM_CLASSIFIER_ENABLED</key>
    <string>true</string>
    <key>HUB_ALARM_LLM_DAILY_LIMIT</key>
    <string>100</string>
    <key>HUB_ALARM_CRITICAL_TYPE_ENABLED</key>
    <string>true</string>
    <key>HUB_ALARM_INTERPRETER_ENABLED</key>
    <string>true</string>
    <key>HUB_ALARM_INTERPRETER_FAIL_OPEN</key>
    <string>true</string>
    <key>HUB_ALARM_INTERPRETER_LLM_DAILY_LIMIT</key>
    <string>200</string>
    <key>HUB_ALARM_ENRICHMENT_ENABLED</key>
    <string>true</string>
    <key>HUB_LLM_RATE_LIMIT_PER_MIN</key>
    <string>120</string>
    <key>HUB_SECRETS_RATE_LIMIT_PER_MIN</key>
    <string>240</string>
    <key>HUB_LLM_MAX_QUEUE</key>
    <string>128</string>
    <key>HUB_LLM_QUEUE_TIMEOUT_MS</key>
    <string>15000</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/hub/hub.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/hub/hub.err.log</string>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
