# bots/hub/launchd/ai.telegram.callback-poller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.telegram.callback-poller</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>scripts/telegram-callback-poller.js</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/alexlee/projects/ai-agent-system/bots/hub</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HUB_AUTH_TOKEN</key>
    <string>__SET_IN_LOCAL_LAUNCHAGENT__</string>
    <key>HUB_CONTROL_CALLBACK_SECRET</key>
    <string>__SET_IN_LOCAL_LAUNCHAGENT__</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/telegram-callback-poller.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/telegram-callback-poller.err</string>
</dict>
</plist>
